# templates/supervisor.program.ini
[program:${process}]
command = ${command}
process_name = ${process}
startsec = 5
redirect_stderr = true
stdout_logfile = ${log_file}
username = ${supervisord.conf:username}
password = ${supervisord.conf:password}
