# templates/supervisor.init
#!/bin/sh
### BEGIN INIT INFO
# Provides:          supervisor instance
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Debian startup script for a supervisor instance
# Description:       This file should be used to construct scripts to be
#                    placed in /etc/init.d.
### END INIT INFO
supervisorctl="${buildout:directory}/bin/supervisorctl"
supervisord="${buildout:directory}/bin/supervisord"
name="${settings:projectname}"
[ -f $supervisord ] || exit 1
[ -f $supervisorctl ] || exit 1
RETVAL=0
start() {
   echo -n "Starting $name: "
   $supervisord  -c ${supervisord.conf:output}
   RETVAL=$?
   [ $RETVAL -eq 0 ] && touch ${supervisord.conf:file}
   echo
   return $RETVAL
}
stop() {
   echo -n "Stopping $name: "
   $supervisorctl -c ${supervisord.conf:output} shutdown
   RETVAL=$?
   [ $RETVAL -eq 0 ] && rm -f ${supervisord.conf:file}
   echo
   return $RETVAL
}
status() {
  $supervisorctl -c ${supervisord.conf:output} status
  RETVAL=$?
  return $RETVAL
}
case "$1" in
       start)
           start
           ;;
       stop)
           stop
           ;;
       status)
           status
           ;;
       *)
           echo "$0 [start|stop|status]"
           ;;
esac
exit $REVAL
